# samba/fixsmb.service
[Unit]
Description=script description %I

[Service]
Type=simple
ExecStart=/usr/bin/fixsmb
Restart=on-failure

[Install]
WantedBy=multi-user.target